# deploy/default.supervisor.conf
[program:zabbix_battery]
command=/path/to/app.py --port=8000
stderr_logfile = /var/log/supervisord/tornado-stderr.log
stdout_logfile = /var/log/supervisord/tornado-stdout.log