# deploy/default.supervisor.conf
[program:zabbix_battery]
command=/usr/bin/python -m zabbix_battery.server
stderr_logfile = /var/log/supervisord/zabbix_battery.error.log
stdout_logfile = /var/log/supervisord/zabbix_battery.log
